# data/systemd/br.inf.ufes.lar.pppoedi.Service.service.in
[Unit]
Description=PPPoEDI Service

[Service]
Type=dbus
BusName=@project_name@.Service
ExecStart=@bindir@/@project_name@.service
StandardError=null
StandardOut=null
Restart=on-abort

[Install]
WantedBy=multi-user.target
